# src/sam_trader/services/crontab
# SAM Trader V3 — Operations cron schedule
# Timezone: Asia/Hong_Kong
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
TZ=Asia/Hong_Kong

# ---------------------------------------------------------------------------
# Daily backup at 05:00 HKT on weekdays (within maintenance window).
# Skips US/HK trading holidays via script logic in backup.py.
# ---------------------------------------------------------------------------
0 5 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.backup backup >> /opt/sam_trader/logs/backup.log 2>&1

# ---------------------------------------------------------------------------
# Log rotation at 03:00 HKT daily.
# ---------------------------------------------------------------------------
0 3 * * * sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.rotate_logs >> /opt/sam_trader/logs/rotate.log 2>&1

# ---------------------------------------------------------------------------
# Deployment window check (configurable via DEPLOY_WINDOW, default 05:00-08:00
# HKT). Runs every 30 minutes during the likely window so custom values are
# honoured without rebuilding the crontab.
# ---------------------------------------------------------------------------
*/30 4-9 * * * sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.deploy_window >> /opt/sam_trader/logs/deploy_window.log 2>&1

# ---------------------------------------------------------------------------
# Pre-market pipeline — HK session
# Runs at 07:30 HKT weekdays (2h before HK market open at 09:30 HKT).
# Gated on MarketCalendarService.is_trading_day('HK') inside pipeline.py.
# ---------------------------------------------------------------------------
30 7 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.pipeline --market HK >> /opt/sam_trader/logs/pipeline_hk.log 2>&1

# ---------------------------------------------------------------------------
# Pre-market pipeline — US session
# Runs at 20:30 HKT weekdays (08:30 ET — 1h before US market open at 09:30 ET).
# Gated on MarketCalendarService.is_trading_day('US') inside pipeline.py.
# ---------------------------------------------------------------------------
30 20 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.pipeline --market US >> /opt/sam_trader/logs/pipeline_us.log 2>&1

# ---------------------------------------------------------------------------
# Performance analysis — nightly stats computation (02:00 HKT daily).
# ---------------------------------------------------------------------------
0 2 * * * sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.performance_analyzer >> /opt/sam_trader/logs/performance.log 2>&1

# ---------------------------------------------------------------------------
# Safety circuit-breaker monitor — runs every minute on weekdays.
# Checks daily PnL, rejection streaks, and connectivity.
# ---------------------------------------------------------------------------
* * * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.safety >> /opt/sam_trader/logs/safety.log 2>&1
