# src/sam_trader/services/crontab
# SAM Trader V3 — Operations cron schedule
# Timezone: Asia/Hong_Kong
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
TZ=Asia/Hong_Kong

# ---------------------------------------------------------------------------
# Daily backup at 06:00 HKT on weekdays; skips US/HK trading holidays via
# script logic in backup.py.
# ---------------------------------------------------------------------------
0 6 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.backup backup >> /opt/sam_trader/logs/backup.log 2>&1

# ---------------------------------------------------------------------------
# Log rotation at 03:00 HKT daily.
# ---------------------------------------------------------------------------
0 3 * * * sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.rotate_logs >> /opt/sam_trader/logs/rotate.log 2>&1

# ---------------------------------------------------------------------------
# Deployment window check (configurable via DEPLOY_WINDOW, default 05:00-08:00
# HKT). Runs every 30 minutes during the likely window so custom values are
# honoured without rebuilding the crontab.
# ---------------------------------------------------------------------------
*/30 4-9 * * * sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.deploy_window >> /opt/sam_trader/logs/deploy_window.log 2>&1

# ---------------------------------------------------------------------------
# Pre-market pipeline — HK session
# Runs at 07:30 HKT weekdays (2h before HK market open at 09:30 HKT).
# ---------------------------------------------------------------------------
30 7 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.pipeline --market HK >> /opt/sam_trader/logs/pipeline_hk.log 2>&1

# ---------------------------------------------------------------------------
# Pre-market pipeline — US session
# Pipeline schedule is derived from market_config.yaml (08:30 ET for US).
# HKT equivalent is computed dynamically via zoneinfo (DST-aware):
#   Summer (EDT): 08:30 ET = 20:30 HKT
#   Winter (EST): 08:30 ET = 21:30 HKT
# When US DST changes (2nd Sun Mar / 1st Sun Nov), shift these entries ±1h.
#
# Pass 1: early gap scan at 04:30 ET  (summer=16:30 HKT, winter=17:30 HKT)
# Pass 2: trended scan at 08:30 ET  (summer=20:30 HKT, winter=21:30 HKT)
# Pass 3: final AI + risk + bundles + report at 09:00 ET  (summer=21:00 HKT, winter=22:00 HKT)
# ---------------------------------------------------------------------------
30 16 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.pipeline --market US --pass 1 >> /opt/sam_trader/logs/pipeline_us.log 2>&1
30 20 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.pipeline --market US --pass 2 >> /opt/sam_trader/logs/pipeline_us.log 2>&1
0 21 * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.pipeline --market US --pass 3 >> /opt/sam_trader/logs/pipeline_us.log 2>&1

# ---------------------------------------------------------------------------
# Performance analysis — nightly stats computation (02:00 HKT daily).
# ---------------------------------------------------------------------------
0 2 * * * sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.performance_analyzer >> /opt/sam_trader/logs/performance.log 2>&1

# ---------------------------------------------------------------------------
# Safety circuit-breaker monitor — runs every minute on weekdays.
# Checks daily PnL, rejection streaks, and connectivity.
# ---------------------------------------------------------------------------
* * * * 1-5 sam . /opt/sam_trader/.env_cron && cd /opt/sam_trader && /usr/local/bin/python3 -m sam_trader.services.safety >> /opt/sam_trader/logs/safety.log 2>&1
